# ansible/templates/etc/systemd/system/redis-server.service.d/permissions.conf.j2
# {{ ansible_managed }}

[Service]
User=pretix
Group=pretix
