# ansible/roles/tavern/templates/deploy-tavern.service.j2
[Unit]
Description=Deploy Tavern

[Service]
Type=oneshot
User=tavern
Group=tavern
ExecStart=/usr/local/bin/deploy-tavern.sh
