# ansible/templates/etc/update-motd.d/01-skyblue-header.j2
#!/bin/sh
# {{ ansible_managed }}

echo ""
toilet -F gay -f future Skyblue
echo ""
