# ansible/templates/etc/update-motd.d/01-skyblue-header.j2
#!/bin/sh

toilet -F gay -f future Skyblue
