# ansible/templates/etc/update-motd.d/02-playbook-revision.j2
#!/bin/sh
# {{ ansible_managed }}

echo "The most recently deployed playbook revision"
echo "is \"{{ playbook_revision | default('unknown') }}\"."
