# ansible/templates/etc/systemd/system/certbot.service.d/override.conf.j2
# {{ ansible_managed }}

[Service]
ExecStart=/usr/bin/certbot renew --quiet --agree-tos --post-hook="/bin/systemctl reload nginx.service"
