# ansible/templates/etc/systemd/system/certbot.timer.d/override.conf.j2
# {{ ansible_managed }}

[Timer]
# Needed to clear any previously existing value for a setting
OnCalendar=
OnCalendar=Sat *-*-* 05:00 UTC
RandomizedDelaySec=
RandomizedDelaySec=900
