# ansible/roles/packages/templates/unattended-upgrades-output.service.j2
# {{ ansible_managed }}

[Unit]
Description=Send output of unattended upgrades log to Slack

[Service]
Type=simple
User=root
Group=root
PrivateTmp=true

{{ lookup('template', 'healthcheck.conf.j2') }}
