# ansible/roles/pretix/templates/systemd/pretix-runperiodic.timer.j2
# {{ ansible_managed }}

[Unit]
Description=Clear expired pretix sessions hourly
After=network.target

[Timer]
OnCalendar=hourly
AccuracySec=5min

[Install]
WantedBy=timers.target
