# ansible/roles/koala/templates/systemd/koala-check-study-year.timer.j2
# {{ ansible_managed }}

[Unit]
Description=Nightly check for creation of new study year for Constipated Koala

[Timer]
OnCalendar=*-*-* 01:00 UTC
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target
