# ansible/templates/etc/systemd/system/backup-admins.service.j2
# {{ ansible_managed }}

[Unit]
Description=Backup of the contents of the admins' home directories to AWS S3

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/backup-admins.sh
