# ansible/templates/etc/systemd/system/db-backup.service.j2
# {{ ansible_managed }}

[Unit]
Description=Database backup to S3

[Service]
Type=simple
User=ansible
ExecStart=/home/ansible/.aws/db-backup.sh
