# ansible/templates/etc/systemd/system/db-backup.timer.j2
# {{ ansible_managed }}

[Unit]
Description=Nightly backup of all SQL databases

[Timer]
OnCalendar=Sat *-*-* 06:00:00
Persistent=true

[Install]
WantedBy=timers.target
