# ansible/templates/etc/systemd/system/db-backup.service.j2
# {{ ansible_managed }}

[Unit]
Description=Backup of all MariaDB databases to AWS S3

[Service]
Type=simple
User=ansible
Group=ansible
ExecStart=/home/ansible/.aws/db-backup.sh
