# ansible/templates/etc/systemd/system/deploy-static-sticky.service.j2
# {{ ansible_managed }}

[Unit]
Description=Deploy static-sticky

[Service]
Type=oneshot
User=static-sticky
Group=static-sticky
ExecStart=/usr/local/bin/deploy-static-sticky.sh
