# ansible/templates/etc/systemd/system/koala-studystatus-reminder.timer.j2
# {{ ansible_managed }}

[Unit]
Description=Weekly Slack reminder sbout running studystatus

[Timer]
OnCalendar=Sat *-*-* 08:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
