# ansible/templates/etc/systemd/system/oauth2-proxy.service.j2
[Unit]
Description=The Koala authentication daemon
[Service]
Restart=always
ExecStart=/usr/local/bin/oauth2_proxy -config /etc/oauth2_proxy/config.cfg
