# ansible/templates/etc/systemd/system/netdata.service.d/onfailure.conf.j2
# {{ ansible_managed }}

[Unit]
OnFailure=failure-notificator@%n.service
