# ansible/templates/etc/systemd/system/netdata.service.d/permissions.conf.j2
# {{ ansible_managed }}
[Service]
Group=www-data

RuntimeDirectory=netdata
RuntimeDirectoryMode=0750
