# ansible/roles/certbot/templates/certbot.service.d/override.conf.j2
# {{ ansible_managed }}

[Service]
{{ lookup('template', 'healthcheck.conf.j2') }}
ExecStartPre=/usr/bin/certbot renew --quiet --agree-tos
