# ansible/templates/etc/systemd/system/mysql.service.d/override.conf.j2
# {{ ansible_managed }}

[Service]
Restart=always
