# ansible/templates/etc/systemd/system/nginx.service.d/umask.conf.j2
# {{ ansible_managed }}

[Service]
UMask=007
