# ansible/roles/intro_website/templates/deploy-intro-website.service.j2
[Unit]
Description=Deploy intro website

[Service]
Type=oneshot
User=intro-website
Group=intro-website
ExecStart=/usr/local/bin/deploy-intro-website.sh
